# crontab
0 0 * * * python /app/app.py >> /var/log/cron.log
